# dist/systemd/aurora-alarm-daemon.service
[Unit]
Description=Aurora Alarm Daemon
After=graphical-session.target

[Service]
Type=dbus
BusName=io.codex.AuroraAlarm
ExecStart=%h/.local/bin/alarm-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
